# reverse_tunnel.conf
description "Reverse SSH tunnel"

start on [2345]
stop on [!2345]

respawn
respawn limit 10 60

exec /usr/local/bin/reverse_tunnel
